# server/wz-server1.service.example
# WastelandZ Game Server 1 — systemd service (Linux)
# 1. Copy to /etc/systemd/system/wz-server1.service   (drop the .example)
# 2. Change "youruser" below to your Linux username
# 3. sudo systemctl daemon-reload
# 4. sudo systemctl enable --now wz-server1
[Unit]
Description=WastelandZ Game Server 1
After=network.target

[Service]
WorkingDirectory=/opt/reforger/server
ExecStart=/opt/reforger/server/ArmaReforgerServer -config /opt/reforger/configs/server1.json -profile /opt/reforger/profiles/server1 -addonsDir /opt/reforger/workshop -maxFPS 60
Restart=on-failure
User=youruser
# Big box running 2+ servers? Uncomment to give this server its own cores
# and a RAM safety cap. Core numbers differ per CPU — check with: lscpu -e
# (full walkthrough: the guide's "Big box" section in Part 3)
#CPUAffinity=0-7 16-23
#MemoryMax=32G

[Install]
WantedBy=multi-user.target
